# systemd/rpi-collector.timer
[Unit]
Description=rpi collector timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=30sec

[Install]
WantedBy=timers.target
